# autostart/parcellite.desktop
[Desktop Entry]
Name=Parcellite
Comment=Clipboard with history
Exec=/usr/bin/parcellite -n
Terminal=false
Type=Application
NoDisplay=true
