# autostart/urxvtd.desktop
[Desktop Entry]
Name=Rxvt-Unicode Daemon
Comment=Start urxvt daemon
Exec=urxvtd -q -o
Terminal=false
Type=Application
NoDisplay=true
