# config/autostart/conky.desktop
[Desktop Entry]
Name=Conky
Comment=Conky system monitor
Exec=conky
Terminal=false
Type=Application
NoDisplay=true