# config/autostart/mpd.desktop
[Desktop Entry]
Name=Music Player Daemon
Comment=Start MPD daemon for current user
Exec=mpd
Terminal=false
Type=Application
NoDisplay=true