# config/autostart/parcellite.desktop
[Desktop Entry]
Name=Parcellite
Comment=Clipboard with history
Exec=/usr/bin/parcellite -n -d
Terminal=false
Type=Application
NoDisplay=true