# config/autostart/tint2.desktop
[Desktop Entry]
Name=Tint2
Comment=System tray
Exec=tint2
Terminal=false
Type=Application
NoDisplay=true
