# config/autostart/xcompmgr.desktop
[Desktop Entry]
Name=xcompmgr
Comment=X Composition Manager
Exec=xcompmgr -cfF -D3
Terminal=false
Type=Application
NoDisplay=true